# scripts/pdf2img/pdf2img.service
[Unit]
Description=PDF to Image microservice for n8n
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/pdf2img
EnvironmentFile=/etc/cathedral/pdf2img.env
ExecStart=/usr/bin/python3 /opt/pdf2img/service.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
